# config/rc.local
#!/bin/sh -e
# rc.local for AutoTill
sudo AutoTill.py
exit 0
